# decnet/web/templates/decnet-updater.service.j2
[Unit]
Description=DECNET self-updater (accepts tarball pushes from master) — {{ agent_name }}
Documentation=https://github.com/anti/DECNET
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/decnet
Environment=DECNET_SYSTEM_LOGS=/var/log/decnet/decnet.updater.log
ExecStart=/usr/local/bin/decnet updater --updater-dir /etc/decnet/updater --install-dir /opt/decnet --agent-dir /etc/decnet/agent
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/decnet/decnet.updater.log
StandardError=append:/var/log/decnet/decnet.updater.log

[Install]
WantedBy=multi-user.target
